# services/xencommons.service
[Unit]
Description=XenCommons Init Script

[Service]
ExecStart=/home/dan/bin/xencommons.sh
